# linux/arkane/arkdep-build.d/awzmlinux/overlay/post_bootstrap/etc/systemd/system/suspend-hyprland.service
[Unit]
Description=Suspend hyprland
Before=systemd-suspend.service
Before=systemd-hibernate.service

[Service]
Type=oneshot
ExecStart=/usr/local/bin/suspend-hyprland suspend

[Install]
WantedBy=systemd-suspend.service
WantedBy=systemd-hibernate.service
